# etc/tb-recv.socket
[Unit]
Description=tiltboard receive socket

[Socket]
ListenStream=/run/tb-recv.sock
SocketUser=nginx

[Install]
WantedBy=sockets.target